# etc/prediction_etl_s1.service
[Unit]
Description=cbp-prediction-etl-s1-pipe
After=network.target

[Service]
WorkingDirectory=/home/armandli/cbp-recorder
ExecStart=/bin/bash -c 'cd /home/armandli/cbp-recorder/ && source ./env/bin/activate && prediction_etl_s1 --exchange etl_exchange_s1 --logfile prediction_etl_s1.log'
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/armandli/cbp-recorder/prediction_etl_s1_io.log
StandardError=inherit
Restart=on-failure

[Install]
WantedBy=default.target
