# etc/prediction_ETHUSD_multi_volatility_peachone.service
[Unit]
Description=cbp-prediction-ETHUSD-multi-volatility-peachone
After=network.target

[Service]
WorkingDirectory=/home/armandli/cbp-recorder
ExecStart=/bin/bash -c 'cd /home/armandli/cbp-recorder/ && source ./env/bin/activate && prediction_torch_s1 --exchange ETHUSD_multi_volatility_peachone --model-key ETHUSD_multi_volatility_peachone --logfile prediction_ETHUSD_multi_volatility_peachone.log'
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/armandli/cbp-recorder/prediction_ETHUSD_multi_volatility_peachone_io.log
StandardError=inherit
Restart=on-failure

[Install]
WantedBy=default.target
