# etc/trade_exchange_0_publisher.service
[Unit]
Description=cbp-trade-0-publisher
After=network.target

[Service]
WorkingDirectory=/home/armandli/cbp-recorder
ExecStart=/bin/bash -c 'cd /home/armandli/cbp-recorder/ && source ./env/bin/activate && trade_publisher --exchange trade_exchange_0 --logfile publisher.log'
Environment=PYTHONUNBUFFERED=1
StandardOutput=append:/home/armandli/cbp-recorder/publisher_io.log
StandardError=inherit
Restart=on-failure

[Install]
WantedBy=default.target
